# config/systemd/trader.service
# ── trader.service ──
# Systemd unit for the autonomous trading runtime.
# Designed for Raspberry Pi OS (Bookworm) / Debian-based systems.
#
# Installation:
#   sudo cp config/systemd/trader.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable trader.service
#   sudo systemctl start trader.service
#
# Logs:
#   journalctl -u trader.service -f
#
# Health check:
#   curl http://localhost:3000/health

[Unit]
Description=trader — Autonomous India Trading Runtime
Documentation=https://github.com/shwetanshu/trader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Run as a non-root service user
User=pi
Group=pi

# Working directory — change this to your install path
WorkingDirectory=/home/pi/trader

# Node.js binary path (adjust for your NodeJS version)
ExecStart=/usr/bin/node --experimental-specifier-resolution=node dist/main.js

# Restart behaviour
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Environment variables
Environment=NODE_ENV=production
Environment=TRADER_PORT=3000
Environment=TRADER_DB_PATH=./data/production.db
Environment=TRADER_SCHEDULER_INTERVAL_MS=60000
Environment=TRADER_LOG_LEVEL=info
Environment=TRADER_MARKET_TIMEZONE=Asia/Kolkata

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/trader/data
CapabilityBoundingSet=
AmbientCapabilities=

# Resource limits
MemoryMax=256M
CPUQuota=50%

# Process lifecycle
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=15

# Standard output goes to journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
