# config/systemd/trader.service
# ── trader.service ──
# Systemd unit for the Hetzner-hosted autonomous trading runtime.
# Uses source entrypoints via tsx because the current repo-wide TypeScript
# build is not globally clean yet.
#
# Installation:
#   sudo cp config/systemd/trader*.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now trader-upstox-notifier.service trader-upstox-mcp-local.service trader.service
#
# Logs:
#   journalctl -u trader.service -f
#
# Health check:
#   curl http://127.0.0.1:3001/health

[Unit]
Description=trader — Autonomous India Trading Runtime
Documentation=https://github.com/shwetanshu/trader
After=network-online.target trader-upstox-notifier.service trader-upstox-mcp-local.service
Wants=network-online.target
Requires=trader-upstox-notifier.service trader-upstox-mcp-local.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/trader
EnvironmentFile=-/root/trader/.env
Environment=NODE_ENV=production
Environment=TRADER_PORT=3001
Environment=TRADER_DB_PATH=./data/trader-upstox.db
Environment=TRADER_SCHEDULER_INTERVAL_MS=60000
Environment=TRADER_LOG_LEVEL=info
Environment=TRADER_MARKET_TIMEZONE=Asia/Kolkata
Environment=TRADER_EXECUTION_OPERATOR_BIND_HOST=127.0.0.1
ExecStart=/usr/bin/node --env-file=/root/trader/.env --import tsx src/main.ts
Restart=always
RestartSec=5
NoNewPrivileges=true
ReadWritePaths=/root/trader/data /root/trader/tmp
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=20
StandardOutput=journal
StandardError=journal
MemoryMax=1G
CPUQuota=150%

[Install]
WantedBy=multi-user.target
