# config/systemd/trader-upstox-notifier.service
[Unit]
Description=trader — Upstox notifier receiver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/trader
EnvironmentFile=-/root/trader/.env
ExecStart=/usr/bin/node --env-file=/root/trader/.env --import tsx src/upstox/notifier-main.ts
Restart=always
RestartSec=5
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ReadWritePaths=/root/trader/tmp /root/trader/data

[Install]
WantedBy=multi-user.target
